# launchd/com.cloudflared-roon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cloudflared-roon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>/Users/monty/.cloudflared/config.yml</string>
        <string>run</string>
        <string>roon-bridge</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/monty/Library/Logs/cloudflared-roon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/monty/Library/Logs/cloudflared-roon.log</string>
</dict>
</plist>
